# service/log/run
#!/bin/sh
exec setuidgid nobody multilog s20000000 n5 /var/log/geodns
