# service/log/run
#!/bin/sh
exec setuidgid nobody multilog t s20000000 n5 /var/log/geodns
